# deploy/supervisord.conf
[supervisord]
nodaemon=true
user=root
loglevel=info
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid

[program:backend]
command=uvicorn app.main:app --host 0.0.0.0 --port %(ENV_BACKEND_PORT)s
directory=/app/backend
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=OPENDRAFT_DATA_DIR="%(ENV_OPENDRAFT_DATA_DIR)s",PORT="%(ENV_BACKEND_PORT)s"

[program:collab]
command=node dist/server.js
directory=/app/collab
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=NODE_ENV="production",DATA_DIR="%(ENV_DATA_DIR)s",PORT="%(ENV_COLLAB_PORT)s",HOST="%(ENV_COLLAB_HOST)s"
